# deploy/pumpkin-selfcheck.service
[Unit]
Description=Pumpkin self-check runner
After=network.target

[Service]
Type=oneshot
User=rossp
WorkingDirectory=/home/rossp/pumpkin-v3
EnvironmentFile=-/etc/pumpkin/pumpkin.env
ExecStart=/usr/bin/python3 -m pumpkin ops selfcheck
Nice=10
